# docs/launchd/com.notes-mcp-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notes-mcp-export</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>notes_mcp.export_notes</string>
        <string>--since-days</string>
        <string>7</string>
        <string>--include-body</string>
        <string>--max-notes</string>
        <string>500</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/ajoyce/git-repos/notes_mcp/src</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/ajoyce/git-repos/notes_mcp</string>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/ajoyce/Library/Logs/notes-mcp/export.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ajoyce/Library/Logs/notes-mcp/export.error.log</string>
</dict>
</plist>
